# battery-guardian.service
[Unit]
Description=Battery Guardian — gestor de carga y salud de la batería
Documentation=file:/opt/battery-guardian/README.md
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/battery-guardian -config /etc/battery-guardian/config.yaml
# Necesita root para escribir conservation_mode y para notificar al usuario.
User=root
Restart=on-failure
RestartSec=5
Nice=10
# NOTA: no se usa ProtectHome/NoNewPrivileges porque el daemon necesita lanzar
# notify-send en la sesión del usuario (socket de D-Bus en /run/user/UID), y ese
# sandboxing lo bloquea.

[Install]
WantedBy=multi-user.target
